# ddrive.conf
[program:ddrive]
user=root
environment=PYTHONPATH=/usr/bin
command=/usr/local/bin/gunicorn ddrive:app -c gunicorn.py
directory=/webserver
startsecs=10

[program:nginx]
command=/usr/sbin/nginx